# mnt-pgsql.mount
# Automatically generated by systemd-fstab-generator

[Unit]
SourcePath=/etc/fstab
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
Before=remote-fs.target
After=network.target azfenv.service

[Mount]
What=//your-azure-storage-accountname.file.core.windows.net/share-file-name
Where=/mnt/pgsql
Type=cifs
Options=credentials=/etc/sysconfig/azfile_env,vers=3.0,dir_mode=0700,file_mode=0600,uid=postgres,gid=postgres

[Install]
WantedBy=multi-user.target
