# pgsql_pitr.timer
[Unit]
Description=Postgres PITR Unit

[Timer]
OnCalendar=daily
Persistent=false
[Install]
WantedBy=timers.target